# packages/mediacenter/kodi/system.d/gui-switch.service
[Unit]
Description=GUI switch service
Before=kodi.service tvip.service
After=var.mount userconfig.service

[Service]
Type=oneshot
ExecStart=-/bin/sh -c 'touch /var/lock/start.tvip'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
